# infra/systemd/user/ai-side-brain-personal-literature-radar.service
[Unit]
Description=AI Side-Brain Personal Literature Radar
Documentation=file:%h/workspace/ai-side-brain/personal/README.md

[Service]
Type=oneshot
WorkingDirectory=%h/workspace/ai-side-brain
ExecStart=%h/workspace/ai-side-brain/scripts/run_personal_literature_radar.sh
