# infra/systemd/user/ai-side-brain-personal-literature-radar-brief.service
[Unit]
Description=AI Side-Brain Personal Literature Radar Brief
Documentation=file:%h/workspace/ai-side-brain/personal/README.md

[Service]
Type=oneshot
WorkingDirectory=%h/workspace/ai-side-brain
ExecStart=%h/workspace/ai-side-brain/scripts/build_personal_literature_radar_brief.sh
