# infra/systemd/user/ai-side-brain-team-literature-radar-cycle.timer
[Unit]
Description=Run AI Side-Brain Team Literature Radar daily cycle

[Timer]
OnCalendar=*-*-* 06:00:00
RandomizedDelaySec=0
Persistent=true
Unit=ai-side-brain-team-literature-radar-cycle.service

[Install]
WantedBy=timers.target
