# infra/systemd/user/ai-side-brain-team-literature-radar.service
[Unit]
Description=AI Side-Brain Team Literature Radar
Documentation=file:%h/workspace/ai-side-brain/team/docs/LITERATURE_RADAR.md

[Service]
Type=oneshot
WorkingDirectory=%h/workspace/ai-side-brain
Environment=RADAR_USE_SAVED_DEFAULTS=1
ExecStart=%h/workspace/ai-side-brain/team/scripts/run_literature_radar.sh
